# systemd/system/nft-firewall-reload.service
[Unit]
Description=reload NFT firewall

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/local/bin/nftfirewall
